# python/crontab
*/5 * * * * /usr/bin/python $scirpt_path
